# supervisord.conf
[include]
files = /etc/supervisor/conf.d/*.conf

[program:app]
command=bash
directory=/project 
autostart=true
autorestart=true
startretries=3
user=root

[supervisord]
nodaemon=true

[supervisorctl]